# Project_1/project1.service
[Unit]
Description=project1 daemon

[Service]
Type=simple
ExecStart=/usr/bin/project1

[Install]
WantedBy=multi-user.target